# systemd/fmn-backend@.service
[Unit]
Description=FMN.consumer backend #%i
After=network.target
Documentation=https://github.com/fedora-infra/fmn.consumer/

[Service]
ExecStart=/usr/lib/python2.7/site-packages/fmn/consumer/backend.py
Type=simple
User=root
Group=root
Restart=on-failure

[Install]
WantedBy=multi-user.target
