# ansible/roles/fmn-dev/files/fmn-digests.service
[Unit]
Description=FMN.consumer digester
After=network.target
Documentation=https://github.com/fedora-infra/fmn.consumer/

[Service]
ExecStart=/home/vagrant/.virtualenvs/python2-fmn/bin/python %h/devel/fmn.consumer/fmn/consumer/digests.py
Type=simple

[Install]
WantedBy=multi-user.target
