# ansible/roles/fmn-dev/files/fmn-backend.service
[Unit]
Description=FMN Delivery Service #%i
After=network.target
Documentation=https://github.com/fedora-infra/fmn/

[Service]
ExecStart=/home/vagrant/.virtualenvs/fmn/bin/twistd -n -l - -y %h/devel/usr/share/fmn/delivery_service.tac
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target
